# docker/ready-to-go/files/rucio.ini
[program:rucio-conveyor-submitter]
command=/bin/rucio-conveyor-submitter  --activities 'User Subscriptions' --sleep-time 10 --group-bulk 20
stdout_logfile=/var/log/rucio/conveyor-submitter.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-conveyor-poller]
command=/bin/rucio-conveyor-poller --sleep-time 10
stdout_logfile=/var/log/rucio/conveyor-poller.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-conveyor-finisher]
command=/bin/rucio-conveyor-finisher --sleep-time 10
stdout_logfile=/var/log/rucio/conveyor-finisher.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-undertaker]
command=/bin/rucio-undertaker
stdout_logfile=/var/log/rucio/undertaker.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-reaper]
command=/bin/rucio-reaper --total-workers 10 --greedy --include-rses reaper=1
stdout_logfile=/var/log/rucio/reaper.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1
environment=
   GLOBUS_THREAD_MODEL=pthread,
   X509_USER_PROXY=%(ENV_X509_USER_PROXY)s,
   X509_USER_KEY=%(ENV_X509_USER_PROXY)s,
   X509_USER_CERT=%(ENV_X509_USER_PROXY)s

[program:rucio-necromancer]
command=/bin/rucio-necromancer
stdout_logfile=/var/log/rucio/necromancer.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-abacus-account]
command=/bin/rucio-abacus-account
stdout_logfile=/var/log/rucio/abacus-account.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-abacus-rse]
command=/bin/rucio-abacus-rse
stdout_logfile=/var/log/rucio/abacus-rse.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-transmogrifier]
command=/bin/rucio-transmogrifier
stdout_logfile=/var/log/rucio/transmogrifier.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-judge-evaluator]
command=/bin/rucio-judge-evaluator
stdout_logfile=/var/log/rucio/judge-evaluator.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-judge-cleaner]
command=/bin/rucio-judge-cleaner
stdout_logfile=/var/log/rucio/judge-cleaner.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-judge-repairer]
command=/bin/rucio-judge-repairer
stdout_logfile=/var/log/rucio/judge-repairer.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-conveyor-receiver]
command=/bin/rucio-conveyor-receiver
stdout_logfile=/var/log/rucio/conveyor-receiver.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-hermes]
command=/bin/rucio-hermes
stdout_logfile=/var/log/rucio/hermes.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1

[program:rucio-panoptes]
command=/opt/rucio/tools/panoptes.py  --inotify /home/garvin
stdout_logfile=/var/log/rucio/panoptes.log
childlogdir=/var/log/rucio/
stdout_logfile_maxbytes=50MB
redirect_stderr=true
stopsignal=KILL
exitcodes=1
